# systemd/audioKeepalive.timer
[Unit]
Description=Run keepalive every 10 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
